# debian/linbpq.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/sbin/linbpq -c /etc -d /opt/oarc/linbpq -l /opt/oarc/linbpq/logs
WorkingDirectory=/opt/oarc/bpq
Restart=always
User=linbpq
Group=linbpq

[Install]
WantedBy=multi-user.target
